# systemd/cg-db-writer-cleanup.service
[Unit]
Description=CG DB-Writer — Retention cleanup (oneshot)
After=postgresql.service
Wants=postgresql.service

[Service]
Type=oneshot
User=cg
Group=cg
WorkingDirectory=/opt/db-writer
ExecStart=/opt/db-writer/venv/bin/python -m src --config /etc/db-writer/config.yml --cleanup
StandardOutput=journal
StandardError=journal
