# systemd/cg-db-writer-cleanup.timer
[Unit]
Description=CG DB-Writer — Retention cleanup timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
